# conf/systemd.service
[Unit]
Description=minecraft
After=network.target

[Service]
WorkingDirectory=__PATH__
User=admin
Group=users
Type=simple
UMask=000
ExecStart=/usr/bin/java -Xmx1024M -Xms1024M -jar server.jar nogui
RestartSec=30
Restart=always

[Install]
WantedBy=multi-user.target
