# conf/systemd.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
Nice=1
KillMode=none
SuccessExitStatus=0 1
ProtectHome=true
ProtectSystem=full
PrivateDevices=true
NoNewPrivileges=true
WorkingDirectory=__FINALPATH__/
UMask=000
ExecStart=/usr/bin/java -Xmx__RAM__M -Xms__RAM__M -jar server.jar nogui
ExecStop=__FINALPATH__/mcrcon/mcrcon -H 127.0.0.1 -P __RCONPORT__ -p __PASSWORD__ stop
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
