# conf/systemd.service
[Unit]
Description=minecraft
After=network.target

[Service]
WorkingDirectory=__PATH__
User=admin
Group=users
Type=simple
UMask=000
ExecStart=/usr/bin/java -Xmx__RAM__M -Xms__RAM__M -jar server.jar nogui
RestartSec=15
Restart=always

[Install]
WantedBy=multi-user.target